# grouperDispatcher/src/main/resources/init.d/grouperDispatcher
##!/bin/bash 
#
#
# chkconfig: 345 82 20
#
#    This shell script takes care of starting and stopping
#       the grouper dispatcher application (grouperDispatcher).
#
# The grouperDispatcher is a java application that reads 
# grouper changelog transactionsfrom an ActiveMQ queue 
# and  writes to multiple ActiveMQ Queues for use by 
# different applications based on a configuration file.

export JAVA_HOME=/usr/java/jdk1.7.0_45/jre/
export PATH=${JAVA_HOME}/bin:${PATH}

# echo "JAVA_HOME:  $JAVA_HOME"

export GROUPER_DISPATCHER_HOME=/opt/grouperDispatcher/
prog=grouperDispatcher

start() {
        check_process_running
    CHECK_RET=$?
    if [ $CHECK_RET -eq 1 ];
    then
           echo "grouperDispatcher already running"
       exit 9
        fi  
    echo "Starting ${GROUPER_DISPATCHER_HOME}/${prog} "  
    [ -f ${GROUPER_DISPATCHER_HOME}conf/grouperDispatcher.properties ] || exit 5
    [ -r ${GROUPER_DISPATCHER_HOME}conf/grouperDispatcher.properties ] || exit 6
    [ -f ${GROUPER_DISPATCHER_HOME}conf/grouperDispatcherConfig.txt ] || exit 7
    [ -r ${GROUPER_DISPATCHER_HOME}conf/grouperDispatcherConfig.txt ] || exit 8
    [ -f ${GROUPER_DISPATCHER_HOME}/grouperDispatcher.jar ] || exit 5
    
    nohup java -jar ${GROUPER_DISPATCHER_HOME}grouperDispatcher.jar >${GROUPER_DISPATCHER_HOME}logs/${prog}.out 2>&1 &
    RETVAL=$?
        echo $RETVAL
    return $RETVAL
}

check_process_running() {
    PROCESS_NUM=`ps -ef | grep "grouperDispatcher.jar" | grep -v "grep" | wc -l` 
        if [ $PROCESS_NUM -eq 1 ];
        then
                return 1
        else
                return 0
        fi
}

status() {
    check_process_running
        CHECK_RET=$?
    if [ $CHECK_RET -eq 0 ];
    then
           echo "grouperDispatcher is not  running"
       exit 9
        fi  
        RESULT=`ps -ef | grep "grouperDispatcher.jar" | grep -v "grep"`
    echo $RESULT
}

stop() {
    check_process_running
        CHECK_RET=$?
    if [ $CHECK_RET -eq 0 ];
    then
           echo "grouperDispatcher is not  running"
       exit 9
        fi  
    echo "Stopping ${GROUPER_DISPATCHER_HOME}/${prog} ... This will take a few seconds. be patient"  
    PID=`ps -ef| grep "grouperDispatcher.jar" | grep -v "grep" | awk '{print $2}'`
    RC=$?
    [ "$RC" != "0" ] && exit 4
        echo -n $"Shutting down $prog: "
    kill $PID
    RETVAL=$?
        COUNTER=0
        while [ $COUNTER -lt 15 ]
    do
          COUNTER=$((COUNTER+1))
          sleep 5
      check_process_running
          CHECK_RET=$?
          if [ $CHECK_RET -eq 0 ];
          then
           echo "grouperDispatcher is not  running"
           break 
          fi
        done
        check_process_running
    if [ $CHECK_RET -eq 0 ];
        then
           echo "grouperDispatcher is not  running"
        fi
    return $RETVAL
}

# See how we were called.
case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status $prog
    ;;
  restart|force-reload)
    stop
    start
    ;;
  *)
    echo $"Usage: $0 {start|stop|status|restart}"
    exit 2
esac